# systemd/ares-omokoda-memory.service
[Unit]
Description=ares-omokoda-memory — Julia memory/compute service (Busy Beaver, NIST, Augury, DePIN, REM fractal planning)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ares/Omo-Koda2/omokoda-memory
ExecStart=/usr/local/bin/julia --project=. server.jl --port 7778
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
