# systemd/ares-omokoda-swarm.service
[Unit]
Description=ares-omokoda-swarm — Elixir OTP supervision tree (coordinator, hive, mesh presence, hive-scale REM cycle)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ares/Omo-Koda2/omokoda-swarm
ExecStart=/usr/bin/mix run --no-halt
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
