# systemd/ares-omokoda-oya.service
[Unit]
Description=ares-omokoda-oya — ỌYA Go flow/rhythm service (Sabbath gating, per-agent primitive cooldowns)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ares/Omo-Koda2/omokoda-go
Environment=OYA_HTTP_PORT=8100
ExecStart=/opt/ares/Omo-Koda2/omokoda-go/oya
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
