# systemd/ares-omokoda-obatala.service
[Unit]
Description=ares-omokoda-obatala — Ọbàtálá symbolic reasoning/ethics engine (consent, privacy, Hermetic evaluation)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ares/Omo-Koda2/omokoda-clojure
Environment=OBATALA_PORT=4002
ExecStart=/usr/local/bin/bb obatala.clj
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
